# rpmdev/systemd/buildrpm.timer
[Unit]
Description=Build all new listed releases into RPM packages

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target